# supervisord.conf
; Docs: <http://supervisord.org/configuration.html>.

; Notes:
;  - Shell expansion ("~" or "$HOME") is not supported.  Environment
;    variables can be expanded using this syntax: "%(ENV_HOME)s".
;  - Quotes around values are not supported, except in the case of
;    the environment= options as shown below.
;  - Comments must have a leading space: "a=b ;comment" not "a=b;comment".
;  - Command will be truncated if it looks like a config file comment, e.g.
;    "command=bash -c 'foo ; bar'" will truncate to "command=bash -c 'foo ".

[unix_http_server]
file=%(here)s/.supervisor.sock   ; the path to the socket file
;chmod=0700                 ; socket file mode (default 0700)
;chown=nobody:nogroup       ; socket file uid:gid owner
;username=user              ; default is no username (open server)
;password=123               ; default is no password (open server)

; Security Warning:
;  The inet HTTP server is intended for use within a trusted environment only.
;  It should only be bound to localhost or only accessible from within an
;  isolated, trusted network.  The inet HTTP server does not support any
;  form of encryption.  The inet HTTP server does not use authentication
;  by default (see the username= and password= options to add authentication).
;  Never expose the inet HTTP server to the public internet.

[inet_http_server]          ; inet (TCP) server disabled by default
port=127.0.0.1:3044         ; ip_address:port specifier, *:port for all iface
;username=user              ; default is no username (open server)
;password=123               ; default is no password (open server)

[supervisord]
logfile=supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
identifier=prestige-supervisor       ; supervisord identifier, default is 'supervisor'

; The rpcinterface:supervisor section must remain in the config file for
; RPC (supervisorctl/web interface) to work.  Additional interfaces may be
; added by defining them in separate [rpcinterface:x] sections.

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

; The supervisorctl section configures how supervisorctl will connect to
; supervisord.  configure it match the settings in either the unix_http_server
; or inet_http_server section.

[supervisorctl]
serverurl=unix://%(here)s/.supervisor.sock ; use a unix:// URL  for a unix socket
;serverurl=http://127.0.0.1:3044 ; use an http:// url to specify an inet socket
;username=chris              ; should be same as in [*_http_server] if set
;password=123                ; should be same as in [*_http_server] if set
;prompt=mysupervisor         ; cmd line prompt (default "supervisor")
;history_file=~/.sc_history  ; use readline history if available

[program:prestige-frontend]
process_name=prestige-frontend
directory=frontend
command=%(here)s/frontend/node_modules/.bin/parcel serve src/index.html --dist-dir dist-serve --port 3040
environment=
	NODE_ENV="development"

[program:prestige-backend]
process_name=prestige-backend
directory=backend
command=%(here)s/venv/bin/python manage.py runserver 127.0.0.1:3041
environment=
	DJANGO_SETTINGS_MODULE="prestige.settings",
	PRESTIGE_DATABASE_URL="sqlite:///db.sqlite3",
	PRESTIGE_ENV="development",
	PRESTIGE_SECRET_KEY="=sez4lhfn@nh86)ylzgl(5k*5kkd+la(dzfsisvdk9ezj2958-",
	PRESTIGE_CORS_ORIGINS="http://localhost:3040",
	PRESTIGE_PROXY_DISALLOW_HOSTS="",
	PYTHONUTF8="1",
	PYTHONUNBUFFERED="1"

[program:prestige-docs]
process_name=prestige-docs
directory=docs
command=%(here)s/venv/bin/mkdocs serve --dev-addr 127.0.0.1:3042
environment=
	PYTHONPATH=.

[program:prestige-reverse-proxy]
process_name=prestige-reverse-proxy
directory=.
command=%(here)s/venv/bin/mitmweb --mode reverse:http://localhost:3042 -s ./mitm_routing.py --listen-port 3045 --web-port 3046 --no-web-open-browser

[group:prestige]
programs=prestige-frontend,prestige-backend,prestige-docs,prestige-reverse-proxy
